# conf/ehdd-unattended-upgrades-after.service
[Unit]
Description=Starts ehdd services

[Service]
Type=oneshot
WorkingDirectory=<INSTALL_ROOT>
ExecStart=<INSTALL_ROOT>/ehdd/startup.sh --no-mount
